# 60-detect.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]|sr[0-9]", RUN+="/usr/local/bin/restore.sh %k"
